# systemd/durian-dashboard.service
[Unit]
Description=Durian Dashboard FastAPI Service
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/durian-dashboard
Environment=MQTT_HOST=sci-iot.ddns.net
Environment=MQTT_PORT=1883
Environment=MQTT_TOPIC=durian_farm1/node_sensor
Environment=MQTT_QOS=1
Environment=DB_PATH=/opt/durian-dashboard/data/durian_dashboard.db
Environment=RETAIN_DAYS=14
Environment=REFRESH_SECONDS=3
ExecStart=/opt/durian-dashboard/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
